# rope_ws/devel/.private/rosbuild/etc/catkin/profile.d/10.rosbuild.sh
# generated from rosbuild/env-hooks/10.rosbuild.sh.em

# env variables in develspace
export ROS_ROOT="/home/lar/ros/rope_ws/src/ros/rosbuild"
