# rope_ws/devel/etc/catkin/profile.d/10.rosbuild.sh
/home/lar/ros/rope_ws/devel/.private/rosbuild/etc/catkin/profile.d/10.rosbuild.sh