# rope_ws/devel/etc/catkin/profile.d/10.rosconsole.sh
/home/lar/ros/rope_ws/devel/.private/rosconsole/etc/catkin/profile.d/10.rosconsole.sh