# rope_ws/devel/etc/catkin/profile.d/15.rosbash.fish
/home/lar/ros/rope_ws/devel/.private/rosbash/etc/catkin/profile.d/15.rosbash.fish