# rope_ws/devel/.private/rosbash/etc/catkin/profile.d/15.rosbash.zsh
# generated from rosbash/env-hooks/15.rosbash.zsh.em

. "/home/lar/ros/rope_ws/src/ros/rosbash/roszsh"
