# rope_ws/devel/.private/rosbash/etc/catkin/profile.d/15.rosbash.tcsh
# generated from rosbash/env-hooks/15.rosbash.tcsh.em

. "/home/lar/ros/rope_ws/src/ros/rosbash/rostcsh"
